# systemd/update-dodns.timer
[Unit]
Description=Update Digital Ocean DNS Timer

[Timer]
WakeSystem=false
OnCalendar=*-*-* 1:00:00

[Install]
WantedBy=timers.target